# deployment/supervisor/reverb.conf
[program:reverb]
process_name=%(program_name)s
# Update this path to your backend project root on the VPS
directory=/var/www/html/backend
command=php artisan reverb:start
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
# Ensure this log file exists and is writable by www-data
stdout_logfile=/var/www/html/backend/storage/logs/reverb.log
stopwaitsecs=3600
