# setup/alerts-poll.service
[Unit]
Description=Laravel Alerts Polling Service
After=network.target

[Service]
# UPDATE THIS PATH TO YOUR PROJECT ROOT
WorkingDirectory=/var/www/test.softdares.com/html/backend

# COMMAND TO START POLLING
ExecStart=/usr/bin/php artisan alerts:poll
Restart=always
RestartSec=5

# UPDATE USER/GROUP IF NEEDED
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
