# setup/reverb.service
[Unit]
Description=Laravel Reverb WebSocket Server
After=network.target

[Service]
# UPDATE THIS PATH TO YOUR PROJECT ROOT
WorkingDirectory=/var/www/test.softdares.com/html/backend

# COMMAND TO START REVERB
ExecStart=/usr/bin/php artisan reverb:start --host=0.0.0.0 --port=8000
Restart=always
RestartSec=5

# UPDATE USER/GROUP IF NEEDED (usually www-data for web servers)
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
